# nzbmegasearch/services/nzbmega/run
#!/bin/bash

exec /sbin/setuser abc python /config/nzbmega/NZBmegasearch/mega2.py
